# conf/supervisord.conf.d/supervisord_nginx.conf
[program:Nginx]
command=nginx -p /opt/qs -c /opt/qs/openresty/nginx/conf/nginx.conf
process_name=Nginx
priority=100
autostart=false
autorestart=true
startretries=10
startsecs=0
stdout_logfile=none
